# build/packer/scripts/proxmox/labs/core/pyxtermjs.service
# General structure created from https://copilot.microsoft.com/shares/SVpE8dfAbZVoWPrHNEicu
[Unit]
Description=pyxtermjs terminal emulator service
After=network.target

[Service]
User=flaskuser
Group=flaskuser
WorkingDirectory=/home/flaskuser
ExecStart=/usr/bin/gunicorn -k eventlet -w 1 pyxtermjs:app --bind 127.0.0.1:5000

Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=15
SyslogIdentifier=CRpyxtermjs

[Install]
WantedBy=multi-user.target