# build/packer/scripts/proxmox/labs/core/wss-proxy.service
[Unit]
Description=Gunicorn instance to serve Flask app
After=network.target

[Service]
User=flaskuser
Group=flaskuser
WorkingDirectory=/home/flaskuser
ExecStart=/usr/bin/gunicorn -k gevent -w 1 app:socketio --bind 127.0.0.1:8000

Restart=always
RestartSec=5
KillSignal=SIGQUIT
TimeoutStopSec=15
SyslogIdentifier=cyberrange

[Install]
WantedBy=multi-user.target
