# build/packer/scripts/proxmox/three-tier/loadbalancer/renew-cert.service
# /etc/systemd/system/renew-cert.service
[Unit]
Description=Run script to renew CA cert

[Service]
Type=oneshot
ExecStart=/usr/bin/step ca renew /etc/ssl/certs/nginx-signed.crt /etc/ssl/private/nginx-signed.key --ca-url https://system22h134.itm.iit.edu --root /root/.step/certs/root_ca.crt --force
