# build/packer/scripts/proxmox/three-tier/loadbalancer/renew-cert.timer
# /etc/systemd/system/renew-cert.timer
[Unit]
Description=Run renew-cert.service at 4 AM on a Tuesday

[Timer]
OnUnitActiveSec=4h
Persistent=true


[Install]
WantedBy=timers.target
