# build/packer/scripts/proxmox/api-server/frontend/renew-cert.timer
# /etc/systemd/system/renew-cert.timer
[Unit]
Description=Run renew-cert.service at 4 AM on a Tuesday

[Timer]
OnCalendar=Tuesday 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
