# build/packer/scripts/proxmox/api-server/frontend/renew-cert.service
# /etc/systemd/system/renew-cert.service
[Unit]
Description=Run my custom script

[Service]
Type=oneshot
ExecStart=step ca renew /home/flaskuser/CAcr.crt /home/flaskuser/CAcr.key --ca-url https://system36.rice.iit.edu --root /root/.step/certs/root_ca.crt --force
